# crudoop/crontab

# Please register this in your crontab #
0 12 * * * curl http://localhost/Phpproject/crudoop/mail.php